# node/cust/happywhale/server/install/happywhale_web.conf
#Put this file into /etc/init on the server and then run
# sudo <command> happywhale_web
# where <command> is one of start, stop, restart, status

description "wildbook node.js server"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    cd /home/ken/Wildbook/node
    mv /var/log/happywhale_web.log /var/log/happywhale_web.log.backup
    #
    # I wanted the most recent version of node on the server because of a change between the
    # version of node that I was using and the version allowed by yum for centos. It caused
    # a bug that I would have had to have two different lines of code depending on the
    # version. Here is how I forced the latest node to be installed on centos
    #    sudo npm cache clean -f
    #    sudo npm install -g n
    #    sudo n stable
    #
    #/usr/bin/node app >> /var/log/happywhale_web.log 2>&1
    /usr/local/n/versions/node/0.12.7/bin/node app >> /var/log/happywhale_web.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
